# config/crontab
# m h  dom mon dow   command
@hourly /home/kuotz/www/current/script/runner -e production  /home/kuotz/www/current/bin/cron.rb 
